# debian/ulogd2.preinst
#!/bin/sh
set -e

# Ensure any ulogd-1.x instances are stopped
if [ -x "/etc/init.d/ulogd" ]; then
	invoke-rc.d ulogd stop || true
fi

#DEBHELPER#

